# stock5x/etc/profile.d/env_setup.sh
export env_adb_device_name_prefix="ingenic"
export env_adb_device_use_diffrent_name="y"
